# target/linux/sunxi/base-files/etc/init.d/zhw_default_network
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
    /sbin/ifconfig eth0:1 169.254.1.123 netmask 255.255.255.0
}
