# target/linux/sunxi/base-files/etc/init.d/zhw_default_network
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
	#check if zhw_default_network is enabled
	local this_enable=$(uci get zhw_default_network.zhw_default_network.this_enable)
	if [ $this_enable != "true" ]; then
		logger "zhw_default_network: is not enabled"
		exit 0
	fi
	logger "zhw_default_network: is enabled"
	/sbin/ifconfig eth0:1 169.254.1.123 netmask 255.255.255.0
}
